# deploy/oci/systemd/proedit-worker.service
[Unit]
Description=ProEdit Worker (Celery)
After=network.target redis-server.service
Wants=redis-server.service

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/opt/proedit/backend
Environment=PYTHONUNBUFFERED=1
Environment=OLLAMA_BASE_URL=http://127.0.0.1:11434
Environment=REDIS_URL=redis://127.0.0.1:6379
EnvironmentFile=-/opt/proedit/backend/.env
ExecStart=/opt/proedit/backend/.venv/bin/python -m celery -A app.celery_app worker --pool=solo --concurrency=1 --without-gossip --without-mingle --without-heartbeat -Q ${CELERY_VIDEO_QUEUE:-video},celery --loglevel=info
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target

